# root/etc/s6-overlay/s6-rc.d/svc-mariadb/run
#!/usr/bin/with-contenv bash

service mariadb start